# scripts/rodent.service
[Unit]
Description=StrataSTOR Rodent Service
After=network.target
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=rodent
Group=rodent
ExecStart=/usr/local/bin/rodent serve
WorkingDirectory=/home/rodent/.rodent
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
Environment=RODENT_CONFIG=/home/rodent/.rodent/rodent.yml
Environment=RODENT_PID_PATH=/home/rodent/.rodent/rodent.pid

# Allow a reasonable amount of time for the service to start
TimeoutStartSec=60

# Security hardening notes:
#
# IMPORTANT: Do not enable PrivateTmp, ReadWritePaths, or similar directives
# that create mount namespaces. ZFS operations require visibility of the host
# mount table. Mount namespace isolation causes:
# - "dataset is busy" errors on destroy
# - Filesystems appear unmounted despite canmount=on
# - Mounts exist in Rodent's namespace but not on host
#
# Rodent uses sudo-based operations (privilege.SudoFileOperations) for
# privileged file access to /etc/samba, /etc/hosts, /etc/krb5.conf, etc.
# Security is enforced via:
# - Code-level allowedPaths validation in privilege.SudoFileOperations
# - /etc/sudoers.d/rodent limiting permitted commands
# - Command validation in command.CommandExecutor
#
# See: https://github.com/stratastor/rodent/issues/54

# ProtectSystem=true
PrivateTmp=false
NoNewPrivileges=false

[Install]
WantedBy=multi-user.target